# Phase 4 · Repos 14–17 Live trading, monitoring & deployment VPS deploy, live dashboard, multi-strategy orchestration/16-vps-deployment-guid/orchestrator.service
# /etc/systemd/system/orchestrator.service
#
# systemd unit: Live Trading Orchestrator (repo 14)
#
# The orchestrator runs the Flask API + background scheduler thread.
# Start this BEFORE the webhook server so /trade is available.
#
# Install:
#   sudo cp systemd/orchestrator.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable orchestrator
#   sudo systemctl start orchestrator

[Unit]
Description=Live Trading Orchestrator (repo 14)
Documentation=https://github.com/yourusername/14-live-trading-orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/home/trader/algo-trading/14-live-trading-orchestrator

EnvironmentFile=/home/trader/algo-trading/14-live-trading-orchestrator/.env

# Gunicorn: 1 worker only — the orchestrator is stateful (in-memory positions)
# Multiple workers would have divergent position state
ExecStart=/home/trader/algo-trading/venv/bin/gunicorn \
    --workers 1 \
    --threads 4 \
    --worker-class sync \
    --bind 127.0.0.1:5001 \
    --timeout 60 \
    --keepalive 5 \
    --log-level info \
    --access-logfile /var/log/algo-trading/orchestrator-access.log \
    --error-logfile  /var/log/algo-trading/orchestrator-error.log \
    --pid /tmp/orchestrator.pid \
    "app:create_app()"

ExecReload=/bin/kill -s HUP $MAINPID

# ⚠️  SIGTERM triggers graceful shutdown: closes open positions before exit
# Give it 60s to complete before force-killing
TimeoutStopSec=60
KillSignal=SIGTERM
KillMode=mixed

Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=120s
StartLimitBurst=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/trader/algo-trading /var/log/algo-trading /tmp

LimitNOFILE=65536
MemoryMax=512M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=orchestrator

[Install]
WantedBy=multi-user.target
